# build/freebsd/ports/pkg-plist
bin/mcsm
etc/mcsm.conf.default
etc/rc.d/mcsm
